# glaido.service
[Unit]
Description=Glaido Voice Transcription
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/glaido/glaido.py
WorkingDirectory=/opt/glaido
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
